# env/.config/systemd/user/backup-todo.service
[Unit]
Description=Backup TODO.md to Google Drive

[Service]
Type=oneshot
ExecStart=%h/.local/scripts/backup-todo
