# env/.config/systemd/user/backup-todo.timer
[Unit]
Description=Hourly backup of TODO.md to Google Drive

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
